# Native/LegacyTransportTests/com.yigit.asciicamera.transport-test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yigit.asciicamera.transport-test</string>
    <key>ProgramArguments</key>
    <array>
        <string>TEST_EXECUTABLE_PATH</string>
    </array>
    <key>MachServices</key>
    <dict>
        <key>com.yigit.asciicamera.transport-test</key>
        <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/ascii-camera-transport-test.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ascii-camera-transport-test.log</string>
</dict>
</plist>
